# cputime_controls/remove_cputime_limit_for_whitelisted_progs.cron
MAILTO=""

*/15 * * * * root /usr/local/sbin/remove_cputime_limit_for_whitelisted_progs
